# minecraft/files/crontab
0 23 * * * /home/minecraft-user/backup.sh
